# Vagrantfiles/hardened_files/debian_customization/rc.local
#!/bin/sh

sudo mount -o remount,noexec /dev/shm
sudo mount -o remount,noexec /dev/
sudo mount -o remount,nodev /run